# deploy/crypto-bot.service
[Unit]
Description=Crypto Bot Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/crypto-bot
EnvironmentFile=/opt/crypto-bot/.env
ExecStart=/usr/bin/python3 /opt/crypto-bot/main.py
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=crypto-bot

[Install]
WantedBy=multi-user.target
